# salutatore/salutatore
#!/bin/sh /etc/rc.common

START=50


start() {
        if [ -e /root/salutatore.sh ];
        then
                /root/salutatore.sh &
        fi
}

stop() {
        killall salutatore.sh
}

